# supervisor/app.conf
[supervisord]
nodaemon=true
user=root

[program:gunicorn]
command=/bin/sh /app/backend/start.sh
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log
environment=PATH="/usr/local/bin:/usr/bin:/bin"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
